# home/.config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -d 5 -c 15 -w 30 -f 0
